# crontab
*/57 * * * * cd /docker/renchon; /usr/bin/docker-compose down && /usr/bin/docker-compose up -d && cd /docker/nginx; /usr/bin/docker-compose down && /usr/bin/docker-compose up -d